# third-party/puppet/cgieoss-ftep/templates/geoserver/initscript.sh.epp
#!/bin/sh
#
# Launch GeoServer as a daemon process.

### BEGIN INIT INFO
# Provides:          geoserver
# Required-Start:    $named $remote_fs $syslog
# Required-Stop:     $named $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop GeoServer
# Description:       GeoServer is an open source software server written in Java that allows users to share and edit geospatial data.
### END INIT INFO

# Configuration files:
#
# /etc/default/geoserver
#   If it exists, sourced at the start of this script. It may perform any
#   sequence of shell commands, like setting relevant environment variables.
#
# The file will be checked for existence before being sourced.

NAME=geoserver
LOCKFILE=/var/lock/subsys/$NAME

##################################################
# See if there's a default configuration file
##################################################
if test -f /etc/default/geoserver ; then
  . /etc/default/geoserver
fi

if test -z "$GEOSERVER_USER" ; then
  GEOSERVER_USER="geoserver"
fi

if test -z "$GEOSERVER_HOME" ; then
  GEOSERVER_HOME="/opt/geoserver"
fi

if test -z "$GEOSERVER_DATA_DIR" ; then
  GEOSERVER_DATA_DIR="${GEOSERVER_HOME}/data_dir"
fi

if test -z "$JAVA_HOME" ; then
  JAVA_HOME="/etc/alternatives/jre"
fi

if test -z "$PORT" ; then
  PORT="9080"
fi

if test -z "$STOPPORT" ; then
  STOPPORT="9079"
fi

case "$1" in
    start)
        ps -ef | grep -v grep | grep "${GEOSERVER_DATA_DIR}" >/dev/null
        RETVAL=$?
        if [ $RETVAL -eq 0 ] ; then
            echo "${NAME} is already running"
        else
            echo "Starting ${NAME}"
            su "$GEOSERVER_USER" -c "touch ${GEOSERVER_DATA_DIR}/serverlog"
            echo `date` >>${GEOSERVER_DATA_DIR}/serverlog
            cd "$GEOSERVER_HOME"
            su "$GEOSERVER_USER" -c "$JAVA_HOME/bin/java -server -DGEOSERVER_DATA_DIR=$GEOSERVER_DATA_DIR -Djava.awt.headless=true -Djetty.port=$PORT -DSTOP.PORT=$STOPPORT -DSTOP.KEY=geoserver -jar start.jar >>$GEOSERVER_DATA_DIR/serverlog 2>&1 &"

            RETVAL=$?
            if [ $RETVAL -eq 0 ]; then
                touch "${LOCKFILE}"
            else
                echo "Failure starting ${NAME}"
            fi
        fi
        ;;
    stop)
        ps -ef | grep -v grep |  grep "${GEOSERVER_DATA_DIR}"  >/dev/null
        RETVAL=$?
        if [ $RETVAL -eq 0 ]; then
            echo  "Stopping ${NAME}"
            echo `date` >>${GEOSERVER_DATA_DIR}/serverlog
            cd "$GEOSERVER_HOME"
            su "$GEOSERVER_USER" -c "$JAVA_HOME/bin/java -DSTOP.PORT=$STOPPORT -DSTOP.KEY=geoserver -jar start.jar --stop >> $GEOSERVER_DATA_DIR/serverlog 2>&1"

            RETVAL=$?
            if [ $RETVAL -eq 0 ]; then
                rm "${LOCKFILE}"
            else
                echo "Failure stopping ${NAME}"
            fi
        else
            echo "${NAME} is not running"
        fi
        ;;
    restart|reload)
        $0 stop
        sleep 5
        $0 start
        ;;
    status)
        ps -ef | grep -v grep |  grep "${GEOSERVER_DATA_DIR}"  > /dev/null
        RETVAL=$?
        if [ $RETVAL -eq 0 ]; then
            echo "$NAME is running."
        else
            echo "$NAME is not running."
        fi
        ;;
    *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 1
esac

exit 0
